# staking/monad/infra/config/monad-validator.service.example
[Unit]
Description=Monad Validator
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=monad
Group=monad
WorkingDirectory=/etc/monad
EnvironmentFile=/etc/monad/validator.env
ExecStart=/usr/local/bin/monad-bft --config /etc/monad/config.toml
Restart=always
RestartSec=5
LimitNOFILE=1048576

[Install]
WantedBy=multi-user.target
