# staking/monad/infra/config/status-server.service.example
[Unit]
Description=Monad Status JSON Server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=monad
Group=monad
WorkingDirectory=/opt/monad-status
ExecStart=/usr/bin/env python3 /opt/monad-status/status_server.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
