# scripts/agenthost/claude-remote-b.service
[Unit]
Description=Claude Code Remote Control - openclaw-wired variant (chimera-openclaw)
After=network-online.target openclaw-gateway.service
Wants=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/home/agents/.local/bin/claude-remote-b-start.sh
ExecStop=/usr/bin/tmux kill-session -t claude-remote-b
Environment=HOME=/home/agents
Environment=PATH=/home/agents/.local/bin:/home/agents/.npm-global/bin:/home/agents/.bun/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
Environment=TMUX_TMPDIR=/tmp

[Install]
WantedBy=default.target
