# scripts/agenthost/eth2-quickstart-20260608.service
[Unit]
Description=Claude Code Remote - eth2-quickstart-20260608 (session manager)
After=network-online.target openclaw-gateway.service
Wants=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/home/agents/.local/bin/eth2-quickstart-20260608-start.sh
ExecStop=/usr/bin/tmux kill-session -t eth2-quickstart-20260608
Environment=HOME=/home/agents
Environment=PATH=/home/agents/.local/bin:/home/agents/.npm-global/bin:/home/agents/.bun/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
Environment=TMUX_TMPDIR=/tmp

[Install]
WantedBy=default.target
