# staking/impl/monad/infra/systemd/monad-validator.service
[Unit]
Description=Monad Validator
After=network-online.target
Wants=network-online.target

[Service]
User=monad
Group=monad
WorkingDirectory=/etc/monad
ExecStart=/usr/local/bin/monad-validator --config /etc/monad/config.toml
Restart=on-failure
RestartSec=5
LimitNOFILE=1048576

# Logs via journald; ensure log rotation if file logging is enabled.

[Install]
WantedBy=multi-user.target
